# nslord/cronhook/crontab
* * * * * /root/cronhook/secure-zones.sh >> /var/log/cron.log 2>&1
